# stow/hypr/.config/systemd/user/arche-battery-monitor.timer
[Unit]
Description=Arche battery monitor timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=30sec
AccuracySec=10sec

[Install]
WantedBy=timers.target
