# system/etc/systemd/system/sshd.service.d/10-arche-oom.conf
# Protect the OpenSSH daemon from the kernel last-resort OOM killer.
# Default OOMScoreAdjust is 0 (eligible like any other process). -500 makes
# sshd one of the last things the kernel kills if everything else has already
# died, so a runaway build that escapes oomd still leaves a working SSH
# entry point for recovery.
#
# This protects the daemon, not the per-login shell: PAM-spawned login
# session-N.scope lives under user-UID.slice and inherits user@.service's
# OOMScoreAdjust=100, so it remains eligible to die with other user
# processes. That's intentional — recovery via re-SSH is the goal, not
# keeping a wedged login alive.

[Service]
OOMScoreAdjust=-500
