# system/etc/systemd/user/dms.service.d/arche.conf
# arche drop-in for the dms user service (/usr/lib/systemd/user/dms.service).
# Linked to /etc/systemd/user/dms.service.d/arche.conf by link_system_all.
# See D032.
#
# DMS_DISABLE_MATUGEN — we drive colors from arche's theming system
#   (theming/themes/active → theming/templates/dms/_emit.sh → custom theme).
#   No wallpaper-derived Material You; the palette is fixed and auditable.
# DMS_DISABLE_POLKIT  — hyprpolkitagent is kept as the auth agent (proven on
#   arche, started from autostart.conf). Two polkit agents would race, so
#   dms's built-in one is off.
[Service]
Environment=DMS_DISABLE_MATUGEN=1
Environment=DMS_DISABLE_POLKIT=1
