# system/etc/systemd/system/docker.service.d/50-arche-memory.conf
# Cap the whole docker engine + all containers at a soft 70% of RAM (~10.5 GiB
# on 15 GiB host). Kernel reclaims containerd/docker pages hard past this,
# preventing a runaway container from wedging the host. No MemoryMax — we
# prefer throttle over SIGKILL restart-loop.

[Service]
MemoryHigh=70%
